# scheduler_plugins/slurm/artifacts/slurm_plugin_cookbook/templates/default/supervisord.conf.erb
[unix_http_server]
file = /var/run/supervisord_pcluster_plugin.sock

[supervisord]
logfile = /var/log/supervisord_pcluster_plugin.log
pidfile = /var/run/supervisord_pcluster_plugin.pid

[supervisorctl]
serverurl = unix:///var/run/supervisord_pcluster_plugin.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

<% case node['pcluster']['node_type'] -%>
<% when 'head' -%>
[program:clustermgtd]
command = <%= node['pcluster']['python_root'] %>/clustermgtd
user = <%= node['plugin']['fleet_mgt_user'] %>
environment = HOME="/home/<%= node['plugin']['fleet_mgt_user'] %>",USER="<%= node['plugin']['fleet_mgt_user'] %>",CONFIG_FILE="<%= node['pcluster']['local_dir'] %>/parallelcluster_clustermgtd.conf",SLURM_BINARIES_DIR="<%= node['slurm']['install_dir'] %>/bin"
redirect_stderr = true
stdout_logfile = /var/log/parallelcluster/clustermgtd

<% when 'compute' -%>
[program:computemgtd]
command = <%= node['pcluster']['python_root'] %>/computemgtd
user = <%= node['plugin']['fleet_mgt_user'] %>
environment = HOME="/home/<%= node['plugin']['fleet_mgt_user'] %>",USER="<%= node['plugin']['fleet_mgt_user'] %>",CONFIG_FILE="<%= node['pcluster']['shared_dir'] %>/parallelcluster_computemgtd.conf",SLURM_BINARIES_DIR="<%= node['slurm']['install_dir'] %>/bin"
redirect_stderr = true
stdout_logfile = /var/log/parallelcluster/computemgtd
<% end -%>
